# k3s/ansible/roles/sapslaj.ipv6_shim/files/ipv6-tcp-shim@.service
[Unit]
Description=IPv6 TCP Shim
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/socat TCP6-LISTEN:%i,fork TCP4:127.0.0.1:%i
