# tohru/ansible/roles/sapslaj.vsdd/files/vsdd.timer
[Unit]
Description=VyOS Shitty Dynamic DNS

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
