# tohru/ansible/roles/sapslaj.photography_rebuild/files/photography-rebuild.timer
[Unit]
Description=photography-rebuild

[Timer]
OnCalendar=daily
RandomizedDelaySec=1800
FixedRandomDelay=true

[Install]
WantedBy=timers.target
