# playboy/ansible/roles/sapslaj.zigbee2mqtt/files/zigbee2mqtt.service
[Unit]
Description=zigbee2mqtt
After=network.target

[Service]
ExecStart=/usr/bin/npm start
WorkingDirectory=/opt/zigbee2mqtt
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
