# aqua/ansible/roles/sapslaj.rclone_jobs/files/rclone-sync@.service
[Unit]
Description=rclone sync %i

[Service]
Type=simple
ExecStart=/usr/local/libexec/rclone-sync-%i.sh
