# aqua/ansible/roles/sapslaj.rclone_jobs/files/rclone-sync@.timer
[Unit]
Description=rclone sync %i
Wants=network-online.target

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
