# tohru/ansible/roles/sapslaj.playboy_backup/files/playboy-backup.timer
[Unit]
Description=Playboy Backup

[Timer]
OnCalendar=hourly
RandomizedDelaySec=1800
FixedRandomDelay=true

[Install]
WantedBy=timers.target
