# common/ansible/roles/sapslaj.docker_standalone/files/docker-housekeeper/docker-housekeeper.timer
[Unit]
Description=Cleanup Docker artifacts to preserve disk space (daily)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
