# tohru/ansible/roles/sapslaj.torrent_auto_pause/files/torrent-auto-pause.timer
[Unit]
Description=Torrent auto-pause for Transmission

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
