# aqua/ansible/files/libvirt/virt-backup-backup.service
[Unit]
Description=virt-backup backup

[Service]
Type=oneshot
ExecStart=/root/.local/bin/virt-backup backup
