# aqua/ansible/files/libvirt/virt-backup-backup.timer
[Unit]
Description=virt-backup backup

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
