# zerotwo/system/roles/sapslaj.ipv6_shim/files/ipv6-udp-shim@.service
[Unit]
Description=IPv6 UDP Shim
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/socat UDP6-LISTEN:%i,fork UDP4:127.0.0.1:%i
